# supervisor/conf.d/worker.conf
[program:worker]
command=/usr/bin/php /usr/share/nginx/operation-service/public/worker.php
autostart=true
autorestart=true
priority=5
